# ublox-recorder/linux/99-ublox-load.rules
ATTRS{idVendor}=="1546", ATTRS{idProduct}=="0507", ACTION=="add", RUN+="/ublox_load.sh"
